# systemd/backupd.timer
[Unit]
Description=BackupD scheduler (checks config every minute)

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
